# contract-assistant/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:streamlit]
command=streamlit run app.py --server.address=0.0.0.0 --server.port=8501 --logger.level=error
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=STREAMLIT_BROWSER_GATHER_USAGE_STATS=false,STREAMLIT_CLIENT_LOGGER_LEVEL=error

[program:nginx]
command=nginx -c /app/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
